# deploy/onprem/co.getoffgridai.sdserver.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- Off Grid image-gen node: sd-server (stable-diffusion.cpp) serving juggernaut-xl-v9 on g3.
     gui-domain LaunchAgent (no sudo). MUST bind 0.0.0.0 so the aggregator on S1 reaches it over
     the LAN at offgrid-g3.local:1234. Install: launchctl bootstrap gui/$(id -u) <this plist>.
     Bootout is async — verify it's gone before re-bootstrap or you get error 5 (EIO). -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>co.getoffgridai.sdserver</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/admin/offgrid-app/Off Grid AI.app/Contents/Resources/bin/sd/sd-server</string>
    <string>-m</string>
    <string>/Users/admin/.offgrid/models/juggernaut-xl-v9-Q4_K.gguf</string>
    <string>--listen-ip</string>
    <string>0.0.0.0</string>
    <string>--listen-port</string>
    <string>1234</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>/Users/admin/sd-server.log</string>
  <key>StandardErrorPath</key><string>/Users/admin/sd-server.log</string>
</dict>
</plist>
